# docker-monitoring/vagrant/src/njmon/custom/root.cron
* 0 * * * /usr/local/bin/njmon -s 300 -c 288 -i vm-njmon -p 8181 -X crc2019
